# deploy/systemd/macchiato-shuiyuan-connector.service.in
[Unit]
Description=macchiatoBot Shuiyuan Discourse connector
After=network.target macchiato-automation.service
Wants=macchiato-automation.service

[Service]
Type=simple
User=__MACCHIATO_USER__
Group=__MACCHIATO_USER__
WorkingDirectory=__MACCHIATO_ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=__MACCHIATO_ROOT__/src
ExecStart=__MACCHIATO_ROOT__/.venv/bin/python -m frontend.shuiyuan_integration.connector
Restart=on-failure
RestartSec=5
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
